# deploy/slack_bot.service
[Unit]
Description=Slack Bot Service with Claude Integration
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/slack_bot
Environment=PATH=/home/ubuntu/slack_bot/venv/bin
EnvironmentFile=/home/ubuntu/slack_bot/.env
ExecStart=/home/ubuntu/slack_bot/venv/bin/python slack_supabase_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
